# templates/service/gitblit.service.chkconfig.erb
#!/bin/bash
# chkconfig: 3 21 91
# description: Starts and Stops gitblit
# Source function library.
. /etc/init.d/functions

# change theses values (default values)
GITBLIT_PATH=/opt/gitblit
GITBLIT_BASE_FOLDER=/opt/gitblit/data
GITBLIT_HTTP_PORT=0
GITBLIT_HTTPS_PORT=8443
GITBLIT_USER="gitblit"
source /etc/sysconfig/gitblit
source ${GITBLIT_PATH}/java-proxy-config.sh
JAVA="java -server -Xmx1024M ${JAVA_PROXY_CONFIG} -Djava.awt.headless=true -jar"

LOCK="/var/lock/subsys/gitblit"
RETVAL=0

pid_of_gitblit() {
  pgrep -u ${GITBLIT_USER} -x -f '^java .*gitblit\.jar .*$'
}
status() {
  pid=$(pid_of_gitblit)
  if [ -n "$pid" ]; then
    echo "gitblit (pid $pid) is running..."
    RET=0
  else
    if [ -f "$LOCK" ]; then
      echo $"gitblit dead but subsys locked"
      RET=2
    else
      echo "gitblit is stopped"
      RET=3
    fi
  fi
  return $RET
}

case "$1" in
  start)
    echo -n $"Starting gitblit server: "
    if [ -f $GITBLIT_PATH/gitblit.jar ]; then
      cd $GITBLIT_PATH
      su $GITBLIT_USER -c "$JAVA $GITBLIT_PATH/gitblit.jar --httpsPort $GITBLIT_HTTPS_PORT --httpPort $GITBLIT_HTTP_PORT --baseFolder $GITBLIT_BASE_FOLDER --dailyLogFile" &
      status > /dev/null
      RETVAL=$?

      cnt=10
      while [ $RETVAL != 0 -a $cnt -gt 0 ]
      do
        status > /dev/null
        RETVAL=$?
        sleep 1
        ((cnt--))
      done

      [ $RETVAL = 0 ] && success $"$STRING" || failure $"$STRING"
      echo
      [ $RETVAL = 0 ] && touch "$LOCK"
    else
      failure $"$STRING"
      echo
      echo $"$GITBLIT_PATH/gitblit.jar not found" 
      RETVAL=5
    fi
  ;;

  stop)
    echo -n $"Stopping gitblit server: "
    pid_of_gitblit > /dev/null
    RETVAL=$?
    if [ $RETVAL = 0 ]; then
      if [ -f $GITBLIT_PATH/gitblit.jar ]; then
        cd $GITBLIT_PATH
        su $GITBLIT_USER -c "$JAVA $GITBLIT_PATH/gitblit.jar --baseFolder $GITBLIT_BASE_FOLDER --stop" > /dev/null

        cnt=10
        while [ $RETVAL = 0 -a $cnt -gt 0 ]
        do
          status > /dev/null
          RETVAL=$?
          sleep 1
          ((cnt--))
        done

        [ $RETVAL = 2 ] && rm -f "$LOCK"
        status > /dev/null
        RETVAL=$?
        if [ $RETVAL = 3 ]; then
          success $"$STRING"
          RETVAL=0
        else
          failure $"$STRING"
          RETVAL=1
        fi
        echo
      else
        failure $"$STRING"
        echo
        echo $"$GITBLIT_PATH/gitblit.jar not found"
        RETVAL=5
      fi
    else
      failure $"$STRING"
      RETVAL=7
      echo
      rm -f "$LOCK"
    fi
  ;;

  status)
    status
    RETVAL=$?
  ;;

  force-reload|restart)
      $0 stop
      sleep 5
      $0 start
  ;;

  *)
    echo $"Usage: /etc/init.d/gitblit {start|stop|status|restart|force-reload}"
    exit 2
  ;;
esac

exit $RETVAL
